# files-templates/daemons/supervisor/scoryst-celery-worker.conf
[program:scoryst-celery-worker]
command=/home/deploy/scoryst/bin/celery-worker-start.sh
user=deploy
stdout_logfile=/home/deploy/scoryst/logs/supervisor-celery-worker.log
redirect_stderr=true
