# files-templates/daemons/supervisor/scoryst.conf
[program:scoryst]
command=/home/deploy/scoryst/bin/gunicorn-start.sh
user=deploy
stdout_logfile=/home/deploy/scoryst/logs/supervisor-gunicorn.log
redirect_stderr=true
